# btc-miner-start.service
[Unit]
Description=Demarre le service miner

[Service]
Type=oneshot
ExecStart=/bin/systemctl start btc-miner.service
