# btc-miner-stop.service
[Unit]
Description=Arrete le service miner

[Service]
Type=oneshot
ExecStart=/bin/systemctl stop btc-miner.service
