# btc-miner-start.timer
[Unit]
Description=Demarre le miner a 23h

[Timer]
OnCalendar=*-*-* 23:00:00
Persistent=true

[Install]
WantedBy=timers.target
